# deploy/supervisor/lexverdict-worker.conf
[program:lexverdict-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/lexverdict/artisan queue:work database --queue=documents,default --sleep=3 --tries=3 --timeout=120 --max-time=3600
directory=/var/www/lexverdict
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/lexverdict/worker.log
stopwaitsecs=130
